# supervisor_conf/php.conf
[program:php]
command = php5-fpm
user = root
autostart = true